# deploy/pdc-redbot-webapp.service
[Unit]
Description=PDC Redbot WebApp (SvelteKit/adapter-node)
After=network.target

[Service]
Type=simple
# Pfade ggf. anpassen:
WorkingDirectory=/opt/pdc/pdc-redbot-webapp
EnvironmentFile=/opt/pdc/pdc-redbot-webapp/.env
ExecStart=/usr/bin/node build
Restart=on-failure
RestartSec=5
User=pdc
Group=pdc

# Härtung
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
